# build/lband/autostart/SDRReceiver.desktop
[Desktop Entry]
Encoding=UTF-8
Version=0.9.4
Type=Application
Name=SDRReceiver
Comment=
Exec=SDRReceiver -s /root/.config/Jontisoft/sdr_25E.ini
OnlyShowIn=XFCE;
RunHook=0
StartupNotify=false
Terminal=false
Hidden=false